# max-mode.service
[Unit]
Description=WSB Snake MAX MODE - Power Hour Trading
Documentation=https://github.com/seanfromthefuteurprivate/Intellibot
After=network.target wsb-snake.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/wsb-snake
Environment="PATH=/root/wsb-snake/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/root/wsb-snake/.env
ExecStart=/root/wsb-snake/venv/bin/python run_max_mode.py
Restart=on-failure
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=max-mode

[Install]
WantedBy=multi-user.target
