# ops/systemd/wsb-ops-audit.service
[Unit]
Description=WSB Snake Audit Agent
After=network.target wsb-snake.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/wsb-snake
EnvironmentFile=/root/wsb-snake/.env
ExecStart=/usr/bin/python3 /root/wsb-snake/ops/audit_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
