# ops/systemd/wsb-ops-deploy.service
[Unit]
Description=WSB Snake Deploy Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/wsb-snake
EnvironmentFile=/root/wsb-snake/.env
ExecStart=/usr/bin/python3 /root/wsb-snake/ops/deploy_agent.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
