# wsb_snake/collectors/wsb-snake-telegram.service
[Unit]
Description=WSB Snake Telegram Screenshot Poller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/wsb-snake
EnvironmentFile=/home/ubuntu/wsb-snake/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/wsb-snake/venv/bin/python /home/ubuntu/wsb-snake/wsb_snake/collectors/telegram_poller.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
